# extra/mydynsd-container.upstart
description "Mydyns container"

start on filesystem and started docker.io
stop on runlevel [!2345]
respawn

script
  /usr/bin/docker start -a <container-id>
end script

